# tg-monitor-bot.service
[Unit]
Description=Telegram Server Monitoring Bot
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/bots/tg-server-monitor
EnvironmentFile=/opt/bots/tg-server-monitor/.env
ExecStart=/opt/bots/tg-server-monitor/tg-monitor-bot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

